# debian/cron.d
5 5 * * * root [ -x /usr/lib/packrat/cron/cleanup ] && /usr/lib/packrat/cron/cleanup
