# debian/cron.d
MAILTO="peter.howe@emc.com"

5 5 * * * root [ -x /usr/local/packrat/cron/cleanup ] && /usr/local/packrat/cron/cleanup

